# deploy/systemd/brave-paws.service
[Unit]
Description=Brave Paws local Tailnet server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=harvey
Group=harvey
WorkingDirectory=/mnt/q/repos/separation-tracker
Environment=BRAVE_PAWS_HOST=127.0.0.1
Environment=BRAVE_PAWS_PORT=4310
Environment=BRAVE_PAWS_PUBLIC_BASE_URL=https://tailnet-host.example.ts.net:7447
Environment="BRAVE_PAWS_DATA_DIR=/mnt/q/fermi/brave-paws/data"
Environment=BRAVE_PAWS_CAMERA_UPSTREAM_BASE_URL=http://127.0.0.1:18888/
ExecStart=/usr/bin/env npm --prefix /mnt/q/repos/separation-tracker run server:start
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
