# deploy/systemd/brave-paws-cd-sync.service
[Unit]
Description=Sync Brave Paws backend to latest released main on QUARK
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=BRAVE_PAWS_REPO_ROOT=/home/harvey/services/separation-tracker-live
Environment=XDG_STATE_HOME=/home/harvey/.local/state
WorkingDirectory=/home/harvey/services/separation-tracker-live
ExecStart=/home/harvey/services/separation-tracker-live/deploy/scripts/sync-brave-paws-backend-to-latest-release.sh
