# deploy/systemd/brave-paws.service
[Unit]
Description=Brave Paws local Tailnet server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=harvey
Group=harvey
WorkingDirectory=/mnt/q/repos/separation-tracker
Environment=BRAVE_PAWS_HOST=127.0.0.1
Environment=BRAVE_PAWS_PORT=4310
Environment=BRAVE_PAWS_PUBLIC_BASE_URL=https://tailnet-host.example.ts.net:7447
Environment="BRAVE_PAWS_DATA_DIR=/mnt/q/fermi/brave-paws/data"
Environment=BRAVE_PAWS_CAMERA_UPSTREAM_BASE_URL=http://127.0.0.1:18888/
Environment=BRAVE_PAWS_CAMERA_CONTROL_PROVIDER=command
Environment="BRAVE_PAWS_CAMERA_CONTROL_LABEL=Picam streaming"
Environment="BRAVE_PAWS_CAMERA_STATUS_COMMAND=/mnt/q/repos/separation-tracker/deploy/scripts/brave-paws-picam-camera-control.sh status"
Environment="BRAVE_PAWS_CAMERA_ENABLE_COMMAND=/mnt/q/repos/separation-tracker/deploy/scripts/brave-paws-picam-camera-control.sh enable"
Environment="BRAVE_PAWS_CAMERA_DISABLE_COMMAND=/mnt/q/repos/separation-tracker/deploy/scripts/brave-paws-picam-camera-control.sh disable"
Environment=BRAVE_PAWS_RECORDING_PROVIDER=command
Environment="BRAVE_PAWS_RECORDING_LABEL=Session recording"
Environment="BRAVE_PAWS_RECORDING_STATUS_COMMAND=/mnt/q/repos/separation-tracker/deploy/scripts/brave-paws-picam-recording-control.sh status"
Environment="BRAVE_PAWS_RECORDING_START_COMMAND=/mnt/q/repos/separation-tracker/deploy/scripts/brave-paws-picam-recording-control.sh start"
Environment="BRAVE_PAWS_RECORDING_STOP_COMMAND=/mnt/q/repos/separation-tracker/deploy/scripts/brave-paws-picam-recording-control.sh stop"
Environment="BRAVE_PAWS_RECORDING_VIDEO_HEIGHT=540"
Environment="BRAVE_PAWS_RECORDING_VIDEO_BITRATE=800k"
Environment="BRAVE_PAWS_RECORDING_VIDEO_MAXRATE=1000k"
Environment="BRAVE_PAWS_RECORDING_VIDEO_BUFSIZE=1600k"
Environment="BRAVE_PAWS_RECORDING_AUDIO_BITRATE=96k"
ExecStart=/usr/bin/env npm --prefix /mnt/q/repos/separation-tracker run server:start
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
