# deploy/systemd/brave-paws-staging-refresh.timer
[Unit]
Description=Poll for Brave Paws QUANTUM staging refreshes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
RandomizedDelaySec=30s
Persistent=true
Unit=brave-paws-staging-refresh.service

[Install]
WantedBy=timers.target
